# deploy/agushare-kiosk.service
[Unit]
Description=A股大屏 - Kiosk 电视全屏显示
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:0
ExecStartPre=/bin/sleep 10
ExecStart=/bin/bash /root/agushare/deploy/kiosk.sh
Restart=always
RestartSec=10
KillMode=process
TimeoutStartSec=120

[Install]
WantedBy=graphical.target
